# deploy/systemd/qtss-web.service.example
[Unit]
Description=QTSS web UI (Vite preview — static assets + proxy to qtss-api)
After=network-online.target
Wants=network-online.target
# Optional ordering when API runs as systemd too (adjust unit name):
# After=qtss-api.service

[Service]
Type=simple
WorkingDirectory=/app/qtss/web
Environment=NODE_ENV=production
# Vite reads `system_config` `api.web_dev_proxy_target` when DATABASE_URL is set (same DB as qtss-api).
Environment=DATABASE_URL=postgres://qtss:CHANGEME@127.0.0.1:5432/qtss
# Optional: QTSS_CONFIG_ENV_OVERRIDES=1 and QTSS_API_PROXY_TARGET=... to bypass DB for the proxy target.
ExecStart=/usr/bin/npm run preview:bind
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
