# deploy/systemd/qtss-web.service.example
[Unit]
Description=QTSS web UI (Vite preview — static assets + proxy to qtss-api)
After=network-online.target
Wants=network-online.target
# Optional ordering when API runs as systemd too (adjust unit name):
# After=qtss-api.service

[Service]
Type=simple
WorkingDirectory=/app/qtss/web
Environment=NODE_ENV=production
# Where `qtss-api` listens (same machine). Used only by Vite preview proxy, not sent to browsers.
Environment=QTSS_API_PROXY_TARGET=http://127.0.0.1:8080
ExecStart=/usr/bin/npm run preview:bind
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
